# com.rajesh.backup-manager.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rajesh.backup-manager</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/rajesh/Desktop/automation/backup-manager/target/release/backup-manager</string>
        <string>--config</string>
        <string>/Users/rajesh/Desktop/automation/backup-manager/config.yaml</string>
        <string>backup</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/rajesh/Desktop/automation/backup-manager/logs/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/rajesh/Desktop/automation/backup-manager/logs/stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
